# supervisor/supervisord_plex.conf
[program:plex]
command=/usr/lib/plexmediaserver/start.sh
user=root
autostart=true
autorestart=true